# systemd/docker-nextcloud-cron.timer
[Unit]
Description=Docker Nextcloud cron

[Timer]
OnCalendar=*-*-* 00/4:10:00

[Install]
WantedBy=timers.target

